# packages/webhost-php-preview/DEBIAN/postinst
#!/bin/bash

# Open pinholes in firewall
if [ -f /etc/network/firewall ]; then
	if [ "`egrep -c '\\$I -p tcp --dport 80 -j ACCEPT' /etc/network/firewall`" != "0" ]; then
		if [ "`egrep -c '\\$I -p tcp --dport 4433 -j ACCEPT   # webhost-php-preview' /etc/network/firewall`" = "0" ]; then
			sed -i '/^\$I -p tcp --dport 80 -j ACCEPT/a\$I -p tcp --dport 4433 -j ACCEPT   # webhost-php-preview\' /etc/network/firewall
		fi
		if [ "`egrep -c '\\$I -p tcp --dport 8080 -j ACCEPT   # webhost-php-preview' /etc/network/firewall`" = "0" ]; then
			sed -i '/^\$I -p tcp --dport 80 -j ACCEPT/a\$I -p tcp --dport 8080 -j ACCEPT   # webhost-php-preview\' /etc/network/firewall
		fi
	fi
	/etc/network/firewall
fi

# Start the preview service
if [ -f "/etc/init.d/webhost-php-preview" ]; then
	/usr/sbin/update-rc.d webhost-php-preview defaults 92
	/etc/init.d/webhost-php-preview start
fi
