# packages/webhost-mgmt-base/etc/udev/rules.d/99-webhost-udev.rules
# WEBHOST UDEV Rules
#
# Modify the timeout value for VMware SCSI devices so that
# in the event of a failover, we don't time out (VMWare best practice).
#

# Redhat systems
#ACTION=="add", BUS=="scsi", SYSFS{vendor}=="VMware, " , SYSFS{model}=="VMware Virtual S",   RUN+="/bin/sh -c 'echo 190 >/sys$DEVPATH/device/timeout'"
# Redhat 6+ systems
#ACTION=="add", SUBSYSTEMS=="scsi", ATTRS{vendor}=="VMware  ", ATTRS{model}=="Virtual disk    ", RUN+="/bin/sh -c 'echo 190 >/sys$DEVPATH/timeout'"

# Debian / SuSE / Ubuntu systems
ACTION=="add", SUBSYSTEMS=="scsi", ATTRS{vendor}=="VMware  " , ATTRS{model}=="Virtual disk    ",   RUN+="/bin/sh -c 'echo 180 >/sys$DEVPATH/device/timeout'"
ACTION=="add", SUBSYSTEMS=="scsi", ATTRS{vendor}=="VMware, " , ATTRS{model}=="VMware Virtual S",   RUN+="/bin/sh -c 'echo 180 >/sys$DEVPATH/device/timeout'"
